# launchd/com.amiwrpremium.macontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent template for macontrol. The `macontrol service install` command
  regenerates this file with the actual binary path and your home directory,
  then loads it via `launchctl bootstrap gui/$UID ...`. You should not normally
  edit this file by hand; use `macontrol service {install,start,stop,status}`.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amiwrpremium.macontrol</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/macontrol</string>
        <string>run</string>
        <string>--log-level=info</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/Library/Logs/macontrol/macontrol.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/Library/Logs/macontrol/macontrol.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
